# doc/rpm-build/egroupware.cron
MAILTO=root
# run eGroupware's async services for all domains
*/5 * * * *	apache	/usr/bin/php -q /usr/share/egroupware/phpgwapi/cron/asyncwrapper.php >/dev/null 2>&1
